# package/controller/files/web.init
#!/bin/sh /etc/rc.common

START=63

start_stop() {
	start-stop-daemon $1 -b -m -p /var/run/wifiController.pid -x /sbin/wifiController
	start-stop-daemon $1 -b -m -p /var/run/etherController.pid -x /sbin/etherController
}

start() {
	start_stop -S
}

reload() {
	return
}

stop() {
	start_stop -K
}
